# service/tbm-drm-user.service
[Unit]
Description=Creating a link file for user to access tbm_drm

[Service]
Type=oneshot
ExecStart=/usr/bin/ln -sf /run/tbm_drm /run/user/%U/
