# service/tbm-drm-user.service
[Unit]
Description=Creating a link file for user to access tbm-drm

[Service]
Type=oneshot
ExecStart=/usr/bin/ln -sf /run/tbm-drm /run/user/%U/
