# service/tbm-drm-user.path
[Unit]
Description=Wait for tbm-drm display daemon socket

[Path]
PathExists=/run/tbm-drm
