# service/tbm-drm.path
[Unit]
Description=Wait for tbm_drm daemon socket

[Path]
PathExists=/run/tbm_drm
